# systemd/radio-movel-sdr.service
[Unit]
Description=Rádio Móvel SDR (somente recepção)
After=systemd-user-sessions.service sound.target
[Service]
Type=simple
User=__USER__
Environment=HOME=/home/__USER__
Environment=QT_QPA_PLATFORM=linuxfb
WorkingDirectory=/opt/radio-movel-sdr
ExecStart=/opt/radio-movel-sdr/venv/bin/python -m app.main
Restart=on-failure
RestartSec=5
[Install]
WantedBy=multi-user.target
