# systemd/radio-movel-sdr-user.service
[Unit]
Description=Rádio Móvel SDR (sessão gráfica do pi)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# Esta é uma unidade de usuário: é instalada em /home/pi/.config/systemd/user
# e, portanto, o processo é executado pelo usuário pi, não pelo root.
Environment=HOME=/home/__USER__
# A aplicação usa X11. Estes valores correspondem ao autologin gráfico padrão
# do Raspberry Pi OS; XDG_RUNTIME_DIR é fornecido pelo gerenciador de sessão.
Environment=QT_QPA_PLATFORM=xcb
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/__USER__/.Xauthority
Environment=XDG_RUNTIME_DIR=%t
WorkingDirectory=/opt/radio-movel-sdr
ExecStart=/opt/radio-movel-sdr/venv/bin/python -m app.main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
